# magic_mirror/rootfs/etc/services.d/magic_mirror/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Magic Mirror from S6"

# Run your program
exec /run.sh
